# scripts/resumeai-report.timer
[Unit]
Description=ResumeAI Daily Report — every day at 08:00 Moscow (05:00 UTC)
Requires=resumeai-report.service

[Timer]
# 05:00 UTC = 08:00 Moscow time
OnCalendar=*-*-* 05:00:00
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
